# install/systemd/rbt-vdom-ami-dtmf.service
[Unit]
Description=RBT virtual domophone: AMI DTMF to HTTP door URLs
Documentation=https://github.com/rosteleset/SmartYard-Server/tree/main/install
After=network-online.target asterisk.service
Wants=asterisk.service

[Service]
Type=simple
WorkingDirectory=/opt/rbt/server
ExecStart=/usr/bin/php /opt/rbt/server/cli/vdom_ami_dtmf_listener.php
RestartSec=10
Restart=always
User=rbt
Group=rbt
# каталог и файлы логов создаёт systemd с владельцем rbt (не нужен ручной chown root)
LogsDirectory=rbt_vdom_ami_dtmf
SyslogIdentifier=rbt_vdom_ami_dtmf
StandardOutput=append:/var/log/rbt_vdom_ami_dtmf/stdout.log
StandardError=append:/var/log/rbt_vdom_ami_dtmf/stderr.log

[Install]
WantedBy=multi-user.target
