# install/systemd/tor-bridge-monitor.timer.example
[Unit]
Description=Периодическая проверка Tor и ротация мостов (RBT)

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
